# emoncms/emoncms_lamp_s6/etc/s6-overlay/s6-rc.d/redis/run
#!/command/execlineb -P
redis-server /etc/redis/redis.conf
