# emoncms/emoncms_lamp_s6/etc/s6-overlay/s6-rc.d/mosquitto/run
#!/command/execlineb -P
mosquitto -c /etc/mosquitto/mosquitto.conf
